# minuspod/optional-transcriber/minuspod-transcriber.container
[Unit]
Description=OpenVINO Model Server (Whisper API)
Wants=podman-user-wait-network-online.service
After=podman-user-wait-network-online.service
Requires=minuspod.service
After=minuspod.service

[Container]
Image=docker.io/openvino/model_server:latest-gpu
ContainerName=minuspod-transcriber
HostName=minuspod-transcriber
Network=minuspod.network
User=5000:5000
GroupAdd=keep-groups
AddDevice=/dev/dri:/dev/dri
Volume=%h/openvino/models:/models:ro
Exec=--config_path /models/config.json --rest_port 8001 --rest_bind_address 0.0.0.0,::
Environment=TZ=America/New_York
AutoUpdate=registry
NoNewPrivileges=yes
DropCapability=ALL
Mount=type=tmpfs,destination=/tmp,tmpfs-size=512M,tmpfs-mode=1777
ReadOnly=true

[Service]
ExecStartPre=/bin/sleep 10
Restart=on-failure
CPUQuota=300%
MemoryMax=4G
MemoryReservation=2G

[Install]
WantedBy=default.target
