# immich/optional/immich-public-proxy.container
[Unit]
Description=Immich Public Proxy
Requires=immich-server.service
After=immich-server.service

[Container]
Image=ghcr.io/alangrainger/immich-public-proxy:latest
ContainerName=immich-public-proxy
HostName=immich-public-proxy
Network=immich.network
PublishPort=3003:3000
AutoUpdate=registry
User=1000:1000
Volume=%h/immich-app/ipp-config.json:/app/config.json:ro
Mount=type=tmpfs,destination=/tmp,tmpfs-size=64M,tmpfs-mode=1777
Mount=type=tmpfs,destination=/var/log,tmpfs-size=32M,tmpfs-mode=1777
Environment=TZ=America/New_York
Environment=PUBLIC_BASE_URL=https://photos.yourdomain.com
Environment=IMMICH_URL=http://immich-webserver:2283
NoNewPrivileges=yes
DropCapability=all
ReadOnly=true

[Service]
ExecStartPre=/bin/sleep 15
Restart=on-failure
CPUQuota=25%
MemoryMax=128M
MemoryLow=64M

[Install]
WantedBy=default.target
