# immich/immich-redis.container
[Unit]
Description=Immich Redis
After=podman-user-wait-network-online.service immich-network.service
Wants=immich-network.service

[Container]
AutoUpdate=registry
Network=immich.network
ContainerName=immich-redis
HostName=immich-redis
Image=docker.io/valkey/valkey:9

Environment=TZ=America/New_York

# Named volume for persistence
Volume=immich_redis_data:/data

# ---------- Hardening (uncomment to lock down) ----------
#ReadOnly=true
#NoNewPrivileges=yes
#DropCapability=all
#AddCapability=CAP_SETUID
#AddCapability=CAP_SETGID

# ---------- Temporary filesystems (required if ReadOnly=true) ----------
#Mount=type=tmpfs,destination=/tmp,tmpfs-size=64M,tmpfs-mode=1777
#Mount=type=tmpfs,destination=/run,tmpfs-size=32M,tmpfs-mode=1777
#Mount=type=tmpfs,destination=/var/log,tmpfs-size=32M,tmpfs-mode=1777

[Service]
Restart=on-failure

[Install]
WantedBy=default.target
