# convertx/convertx.container
[Unit]
Description=ConvertX File Converter Service
After=podman-user-wait-network-online.service
Wants=podman-user-wait-network-online.service

[Container]
Image=ghcr.io/c4illin/convertx:latest
ContainerName=convertx
HostName=convertx
AutoUpdate=registry
PublishPort=3002:3000
Volume=%h/convertx/data:/app/data
Environment=TZ=America/New_York
Environment=HTTP_ALLOWED=true
EnvironmentFile=%h/.config/containers/systemd/convertx.env
LogDriver=k8s-file
LogOpt=max-size=10m
LogOpt=max-file=3
NoNewPrivileges=true
ReadOnly=true
DropCapability=all
AddCapability=CAP_NET_BIND_SERVICE
Mount=type=tmpfs,destination=/var/log,tmpfs-size=64M,tmpfs-mode=1777
Mount=type=tmpfs,destination=/tmp,tmpfs-size=512M,tmpfs-mode=1777
Mount=type=tmpfs,destination=/run,tmpfs-size=16M,tmpfs-mode=1777
PodmanArgs=--memory=1g

[Service]
Restart=on-failure

[Install]
WantedBy=default.target
